# resources/systemd/n2klocd.service
[Unit]
Description=n2klocd daemon

[Service]
User=root
ExecStart=/usr/bin/n2kafka /etc/n2klocd/config.json

[Install]
WantedBy=multi-user.target
